# contrib/init/sandrocoind.conf
description "Sandrocoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env SANDROCOIND_BIN="/usr/bin/sandrocoind"
env SANDROCOIND_USER="sandrocoin"
env SANDROCOIND_GROUP="sandrocoin"
env SANDROCOIND_PIDDIR="/var/run/sandrocoind"
# upstart can't handle variables constructed with other variables
env SANDROCOIND_PIDFILE="/var/run/sandrocoind/sandrocoind.pid"
env SANDROCOIND_CONFIGFILE="/etc/sandrocoin/sandrocoin.conf"
env SANDROCOIND_DATADIR="/var/lib/sandrocoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # sandrocoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$SANDROCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run sandrocoind."
        echo "The setting must appear in $SANDROCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Sandrocoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$SANDROCOIND_PIDDIR"
    chmod 0755 "$SANDROCOIND_PIDDIR"
    chown $SANDROCOIND_USER:$SANDROCOIND_GROUP "$SANDROCOIND_PIDDIR"
    chown $SANDROCOIND_USER:$SANDROCOIND_GROUP "$SANDROCOIND_CONFIGFILE"
    chmod 0660 "$SANDROCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$SANDROCOIND_PIDFILE" \
    --chuid $SANDROCOIND_USER:$SANDROCOIND_GROUP \
    --exec "$SANDROCOIND_BIN" \
    -- \
    -pid="$SANDROCOIND_PIDFILE" \
    -conf="$SANDROCOIND_CONFIGFILE" \
    -datadir="$SANDROCOIND_DATADIR" \
    -disablewallet \
    -daemon

